# projects/slackware-mirror/slackware-mirror.cron
#
# Regular cron jobs for slackware and SlackBuilds stuff
#
0 4	* * *	/usr/bin/slackware-mirror.sh 
