# projects/slackware-mirror/slackware-mirror.cron
#!/bin/sh

/usr/sbin/slackware-mirror.sh >> /var/log/slackware-mirror.log
